# openv-homa/rootfs/etc/s6-overlay/s6-rc.d/vclient_pub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the vclient_pub service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "vclient_pub service starting ..."

declare -r runtime_dir="/run/vcontrold"
declare -r file_mqtt_commands="${runtime_dir}/1_mqtt_commands.txt"
declare -r file_mqtt_template="${runtime_dir}/2_mqtt.tmpl"
declare -r file_mqtt_publish="/dev/shm/3_vcontrold_mqtt_pub.sh"

# Fetch MQTT settings
source get_mqtt_settings.sh

# Fetch Vcontrold settings
source get_vcontrold_settings.sh

# Ensure required variables are set
if [[ -z "$REFRESH_RATE" || -z "$VCONTROL_HOST" || -z "$VCONTROL_PORT" ]]; then
    bashio::exit.nok "vclient_pub run: Missing essential environment variables. Check MQTT and Vcontrold settings."
fi

# Do not exit on error
set +e

# Run your program
while true; do
    if vclient -h "$VCONTROL_HOST" -p "$VCONTROL_PORT" -f ${file_mqtt_commands} \
               -t ${file_mqtt_template} > ${file_mqtt_publish}; then
        source ${file_mqtt_publish}
        exitcode=$?
        if [[ $exitcode -ne 0 ]]; then
            bashio::log.error "vclient_pub run: ${file_mqtt_publish} failed with exit code $exitcode."
        fi
    else
        exitcode=$?
        bashio::log.error "vclient_pub run: vclient failed with exit code $exitcode."
    fi
    bashio::log.info "vclient_pub run: Sleeping for $REFRESH_RATE seconds."
    sleep "$REFRESH_RATE"
done
