# openv-homa/rootfs/etc/s6-overlay/s6-rc.d/vclient_sub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the MQTT-based vclient_sub service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "vclient_sub service starting ..."

# Fetch MQTT settings
source get_mqtt_settings.sh

# Fetch Vcontrold settings
source get_vcontrold_settings.sh

# Ensure required environment variables are set
if [[ -z "$MQTT_HOST" || -z "$MQTT_PORT" || -z "$MQTT_USER" || -z "$MQTT_PASSWORD" || -z "$VCONTROL_HOST" || -z "$VCONTROL_PORT" ]]; then
    bashio::exit.nok "vclient_sub run: Missing essential environment variables. Check MQTT and Vcontrold settings."
fi

# Get configuration options from user
config_systemid=$(bashio::config 'homa_system_id')

# Set SETCMDS_FILE and check if file exists
SETCMDS_FILE="/etc/vcontrold/3_mqtt_setcmds.txt"
if [ ! -f "$SETCMDS_FILE" ]; then
    bashio::log.info "vclient_sub run: No setcmds (${SETCMDS_FILE}) found ..."
    tail -f /dev/null  # Wait here indefinitely to prevent the script from exiting
fi

while IFS=$'\t' read -r topic cmd setcmd type; do
    # Skip empty lines
    [ -z "$topic" ] && continue

    # Subscribe to each entry in the file
    bashio::log.info "vclient_sub run: Starting MQTT subscription to topic: ${topic}/on"
    mosquitto_sub -h "$MQTT_HOST" -p "$MQTT_PORT" -u "$MQTT_USER" -P "$MQTT_PASSWORD" -t "${topic}/on" | while read -r payload; do
        # Here is the callback to execute whenever you receive a message
        bashio::log.info "vclient_sub run: Received topic: ${topic}/on"
        bashio::log.info "vclient_sub run: vclient command ${setcmd}, set value: ${payload}"
        vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -o /dev/null -c "${setcmd} ${payload}"
        exitcode=$?
        if [[ $exitcode -ne 0 ]]; then
            bashio::log.error "vclient_sub run: vclient command ${setcmd} failed with exit code ${exitcode}."
        fi
        result=$(vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -c "${cmd}" -J)
        if [[ $type == "FLOAT" ]]; then
            payload=$(jq -r '.[0].value' <<< "$result")
        elif [[ $type == "FLOAT1" ]]; then
            payload=$(printf "%.1f" "$(jq -r '.[0].value' <<< "$result")")
        elif [[ $type == "INT" ]]; then
            payload=$(printf "%.0f" "$(jq -r '.[0].value' <<< "$result")")
        elif [[ $type == "STRING" ]]; then
            payload=$(jq -r '.[0].raw' <<< "$result")
        else
            bashio::exit.nok "Unsupported type '${type}' for command ${cmd}. Supported types are FLOAT, FLOAT1, INT, STRING."
        fi
        bashio::log.info "vclient_sub run: vclient command ${cmd}, returned value: ${payload}"
        bashio::log.info "vclient_sub run: Publishing to MQTT topic: ${topic}"
        mosquitto_pub -h "$MQTT_HOST" -p "$MQTT_PORT" -u "$MQTT_USER" -P "$MQTT_PASSWORD" -t "${topic}" -m "${payload}"
    done &

done < "$SETCMDS_FILE"

# Let the background processes run
wait

bashio::log.info "vclient_sub run: wait exited ..."
tail -f /dev/null  # Wait here indefinitely to prevent the script from exiting
